# kubelet.service
[Unit]
Description=Kubelet service

[Service]
ExecStart=/usr/local/bin/kubelet

[Install]
WantedBy=default.target
